# docs/npid-api-layer/com.npid.fastapi.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.npid.fastapi</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/singleton23/Raycast/prospect-pipeline/docs/npid-api-layer/start-server.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/singleton23/Raycast/prospect-pipeline/docs/npid-api-layer</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardErrorPath</key>
    <string>/tmp/npid-fastapi.err</string>

    <key>StandardOutPath</key>
    <string>/tmp/npid-fastapi.out</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/singleton23/Raycast/prospect-pipeline/docs/npid-api-layer/venv/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
