# scripts/com.singleton23.prospect-pipeline.ended-meeting-outcome-watch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.singleton23.prospect-pipeline.ended-meeting-outcome-watch</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string><REPO_ROOT>/scripts/run-ended-meeting-outcome-watch.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string><REPO_ROOT></string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string><NODE_BIN_DIR>:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>TZ</key>
    <string>America/New_York</string>
    <key>API_BASE</key>
    <string>http://127.0.0.1:8000/api/v1</string>
    <key>WINDOW_DAYS</key>
    <string>7</string>
    <key>LIMIT</key>
    <string>50</string>
  </dict>

  <key>StartInterval</key>
  <integer>1800</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string><LOCAL_LOG_DIR>/ended-meeting-outcome-watch/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string><LOCAL_LOG_DIR>/ended-meeting-outcome-watch/launchd.err.log</string>
</dict>
</plist>
